# packaging/statefs-system.service.in
[Unit]
Description=StateFS FUSE filesystem, system-wide
# statefs has not this dependency but some providers has
# TODO: make it smarter, maybe it is possible to make
# this dependency conditional
Requires=dbus.socket

# Lets make statefs start as early as possible.
DefaultDependencies=no
After=local-fs.target dbus.socket
Before=basic.target
Conflicts=shutdown.target

[Service]
ExecStartPre=@prefix@/bin/statefs-prerun
ExecStart=@prefix@/bin/statefs /run/state -f --system -o allow_other,gid=@USERS_GID@
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target actdead-pre.target
